# units/admin.target
[Unit]
Description=Admin host container
